# docker-memory-limit.service
[Unit]
Description=docker-memory-limit

[Service]
ExecStart=/root/docker-memory-limit.sh
Restart=always

[Install]
WantedBy=multi-user.target
